# com.maude-analyzer.weekly-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maude-analyzer.weekly-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alexbenison/maude-analyzer/scripts/cron_weekly_refresh.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/alexbenison/maude-analyzer</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/alexbenison/maude-analyzer/logs/weekly_refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alexbenison/maude-analyzer/logs/weekly_refresh_error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
